# system/etc/systemd/user/auto-brightness.service
[Unit]
Description=Self-learning ambient brightness controller
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/auto-brightness.py
Restart=always
RestartSec=5

[Install]
WantedBy=graphical-session.target
