# system/etc/systemd/user/gdm-wallpaper-sync.service
[Unit]
Description=Sync user wallpaper to GDM greeter
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/gdm-wallpaper-sync.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
