# system/etc/systemd/system/anydesk.service.d/restart.conf
# anydesk periyodik SIGABRT ile cokuyordu ve varsayilan unit'te Restart yok ->
# systemd otomatik baslatmiyordu (manuel restart gerekiyordu). Sonsuz self-heal:
# her cokme/exit'te 5sn sonra tekrar baslat, start-limit'e takilma (asla pes etme).
[Unit]
StartLimitIntervalSec=0

[Service]
Restart=always
RestartSec=5
