# system/etc/systemd/system/anydesk.service.d/wait-online.conf
# Boot'ta anydesk.service ag baglanmadan once basliyordu: relay'e ulasamayip
# timeout'larla ugrasarak dakikalarca "AnyDesk not connected" gosteriyordu.
# Cozum: servisi network-online.target'tan sonra siraya al ve baslamadan once
# GERCEK internet erisimini bekle. Olculen: hazir network'te relay ~10s'de
# baglaniyor; nm-online tek basina NM "connected" dedigi an donuyor ama gercek
# internet henuz yokken anydesk relay TCP timeout -> ~2dk backoff. Bu yuzden
# nm-online'dan sonra connectivity GERCEKTEN "full" olana kadar (aktif prob,
# ~60s cap) bekliyoruz; sonra anydesk basliyor ve ilk denemede baglaniyor.
# (Global boot'u yavaslatmamak icin NetworkManager-wait-online'i acmak yerine
#  bekleme sadece bu servise ozel; arka plan servisi oldugu icin login'i tutmaz.)
[Unit]
Wants=network-online.target
After=network-online.target

[Service]
ExecStartPre=/usr/bin/nm-online -t 30 -q
ExecStartPre=/bin/sh -c 'i=0; while [ "$i" -lt 30 ]; do [ "$(nmcli networking connectivity check 2>/dev/null)" = "full" ] && exit 0; i=$((i+1)); sleep 2; done'
