# system/etc/systemd/system/anydesk.service.d/wait-online.conf
# Boot'ta anydesk.service ag baglanmadan once basliyordu: relay'e ulasamayip
# timeout'larla ugrasarak dakikalarca "AnyDesk not connected" gosteriyordu.
# Cozum: servisi network-online.target'tan sonra siraya al ve baslamadan once
# nm-online ile gercek baglantiyi bekle. Boylece ilk relay denemesi basarili olur.
# (Global boot'u yavaslatmamak icin NetworkManager-wait-online'i acmak yerine
#  bekleme sadece bu servise ozel.)
[Unit]
Wants=network-online.target
After=network-online.target

[Service]
ExecStartPre=/usr/bin/nm-online -t 30 -q
