# system/etc/systemd/system/powertop.service.d/override.conf
[Service]
ExecStartPost=/usr/bin/udevadm trigger --subsystem-match=usb --action=change
