# system/etc/systemd/system/cachyos-rate-mirrors.service.d/wait-online.conf
# Boot'ta (Persistent timer catch-up) network-online'a ragmen DNS daha hazir
# degilken calisip "Could not resolve host: geoip.kde.org" ile fail oluyordu.
# Gercek connectivity + DNS cozumu hazir olana kadar bekle (bounded ~60s).
[Service]
ExecStartPre=/bin/sh -c 'i=0; while [ "$i" -lt 30 ]; do [ "$(nmcli networking connectivity check 2>/dev/null)" = "full" ] && getent hosts geoip.kde.org >/dev/null 2>&1 && break; i=$((i+1)); sleep 2; done; exit 0'
