# system/etc/systemd/system/anydesk-watchdog.service
[Unit]
Description=AnyDesk relay watchdog (restart service if relay is not online)
After=network-online.target anydesk.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/anydesk-watchdog
