# user/.config/systemd/user/audio-autoswitch.service
[Unit]
Description=Auto-switch EasyEffects output to the last-plugged audio device
After=pipewire.service wireplumber.service
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/audio-autoswitch
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
